# mqttpi-bms.service
# mqttpi JBD BMS bridge — assumes repo cloned to ~/mqttpi for User=pi
# Edit WorkingDirectory and -c/-s paths if you installed elsewhere.

[Unit]
Description=mqttpi JBD BMS MQTT bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=%h/mqttpi
ExecStart=/usr/bin/python3 -m mqttpi.bms.bridge -c %h/mqttpi/config.yaml -s %h/mqttpi/secrets.yaml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target